# debian/talis-agent.service
[Unit]
Description=Talis Agent - System monitoring agent
After=network.target
Documentation=https://github.com/celestiaorg/talis-agent

[Service]
Type=simple
User=talis-agent
Group=talis-agent
ExecStart=/usr/bin/talis-agent -config /etc/talis-agent/config.yaml
Restart=always
RestartSec=10
WorkingDirectory=/var/lib/talis-agent

# Security settings
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelModules=yes

[Install]
WantedBy=multi-user.target 